# deploy/codex-quota-monitor.timer
[Unit]
Description=Run Codex quota monitor every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
AccuracySec=1min
Unit=codex-quota-monitor.service

[Install]
WantedBy=timers.target
